# crontab
* * * * * root env - `cat /root/env.txt` php /var/www/app/artisan schedule:run

