# crontab
* * * * * root env - `cat /root/env.txt` php /data/www/artisan schedule:run >> /var/log/cron.log &2>1

